# create-install-image/templates/installimg/8.3.0/etc/systemd/system/installer.service
[Unit]
Description=XCP-ng Installer
Requires=systemd-udev-settle.service interface-rename-sideway.service early-blacklist.service
After=systemd-udev-settle.service interface-rename-sideway.service early-blacklist.service getty@tty2.service

[Service]
Type=forking
SendSIGHUP=yes
ExecStart=/opt/xensource/installer/preinit --

[Install]
WantedBy=multi-user.target
