# create-install-image/templates/installimg/8.3.0/etc/systemd/system/sshd.service.d/sshd_installer.conf
[Unit]
ConditionKernelCommandLine=sshpassword

[Service]
ExecStartPre=/bin/sh -c '/bin/sed -e "s/.*sshpassword=\([^ ]*\).*/root:\1/" /proc/cmdline | /sbin/chpasswd'
